# docker/musicbrainz-server-cron/crontab
0 * * * * /home/musicbrainz/musicbrainz-server/admin/cron/hourly.sh
10 0 * * * /home/musicbrainz/musicbrainz-server/admin/cron/daily.sh
